# debian_8/etc/rc.local
#!/bin/sh -e

echo 600 > /sys/block/sda/device/timeout
echo 600 > /sys/block/sr0/device/timeout

exit 0
